# check.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/check.sh

[Install]
WantedBy=default.target


